# packages/tinkerelec/retroarch/system.d/tmp-emulation-retroarch-cores.mount
[Unit]
Description=Cores directory
After=systemd-tmpfiles-setup.service
DefaultDependencies=no

[Mount]
What=none
Where=/tmp/emulation/retroarch/cores
Type=overlay
Options=lowerdir=/usr/lib/libretro,upperdir=/storage/.config/retroarch/cores,workdir=/storage/.cache/retroarch/cores-workdir

[Install]
WantedBy=kodi.target
