# packages/tinkerelec/addons/script.launcher.retroarch/source/system.d/retroarch.service
[Unit]
Description=RetroArch is the reference frontend for the libretro API
After=graphical.target
Requires=graphical.target

[Service]
ExecStart=/bin/sh /storage/.kodi/addons/script.launcher.retroarch/bin/retroarch.start
TimeoutStopSec=30
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
